# debian/prerm
#!/bin/bash

rm /etc/nginx/sites-enabled/radicale.nginx

systemctl restart nginx

systemctl stop radicale
systemctl disable radicale

python3 -m pip uninstall -y radicale

